# time-sync/crontab
41 21 * * * /bin/echo "test" >> /var/log/cron.log 2>&1
# This extra line makes it a valid cron